# systemd/bird-away.service
[Unit]
Description=Bird-Away pool bird deterrent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=gpio
WorkingDirectory=/home/pi/bird-away
EnvironmentFile=/home/pi/bird-away/.env
ExecStart=/home/pi/bird-away/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
